# captcha/etc/service/nginx/run
#!/bin/sh
set -e
exec /usr/sbin/nginx -c /etc/nginx/nginx.conf -g 'daemon off;' >> /var/log/captcha/nginx.access.log 2>&1
